# deployment/self-host/systemd/rowset-backup.service
[Unit]
Description=Create and verify a coordinated Rowset backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/rowset-backup.conf
UMask=0077
TimeoutStartSec=6h
ExecStart=/bin/sh -c 'exec /usr/sbin/runuser -u "$ROWSET_BACKUP_USER" -- "$ROWSET_ROOT/deployment/self-host/backup.sh" "$ROWSET_BACKUP_DIR" "$ROWSET_ENV_FILE"'
